# crontab
0 8,16,0 * * * cd /app && rake synchronize_properties:trovit >> /var/log/cron.log 2>&1
